# Encrypted/arch-encrypted-main/systemd/btrfs85.timer
[Unit]
Description=Cleanup BTRFS once a week
Documentation=man:btrfs

[Timer]
OnCalendar=monthly
AccuracySec=dh

[Install]
WantedBy=timers.target
