# Encrypted/arch-encrypted-main/systemd/btrfs50.timer
[Unit]
Description=Cleanup BTRFS once a week
Documentation=man:btrfs

[Timer]
OnCalendar=weekly
AccuracySec=1h

[Install]
WantedBy=timers.target
