# Encrypted/arch-encrypted-main/systemd/system-upgrade.service
[Unit]
Description=Upgrade system

[Service]
Type=oneshot
ExecStart=/root/bin/system-upgrade
